# Chapter13/deploy/supervisor_worker.conf
[program:celery]
command=celery worker -A celery_runner
directory=/srv/app

